# etc/systemd/system/syncrepo-manjaro.timer
[Unit]
Description=Hourly Manjaro repository sync

[Timer]
OnCalendar=hourly
AccuracySec=1m
Persistent=true

[Install]
WantedBy=timers.target
